# dist/kanku-scheduler.service
[Unit]
Description=Kanku Scheduler
#Documentation=http://wiki2.dovecot.org/
After=local-fs.target network.target
 
[Service]
Type=forking
PIDFile=/opt/kanku/var/run/kanku-scheduler.pid
ExecStart=/opt/kanku/bin/kanku-scheduler

[Install]
WantedBy=multi-user.target
