# dist/kanku-worker.service
[Unit]
Description=Kanku Worker
#Documentation=man:dovecot(1)
#Documentation=http://wiki2.dovecot.org/
After=local-fs.target network.target
 
[Service]
Type=simple
ExecStart=/opt/kanku/bin/kanku-worker
#ExecReload=/usr/sbin/dovecot reload
NonBlocking=yes

[Install]
WantedBy=multi-user.target
