# dist/systemd/kanku-dispatcher.service
[Unit]
Description=Kanku Dispatcher
After=local-fs.target network.target
Wants=rabbitmq-server
 
[Service]
User=kankurun
Type=forking
PIDFile=/run/kanku/kanku-dispatcher.pid
ExecStart=/usr/sbin/kanku-dispatcher
ExecStop=/usr/sbin/kanku-dispatcher --stop

[Install]
WantedBy=multi-user.target
