# dist/systemd/kanku-triggerd.service
[Unit]
Description=Kanku TriggerD
#Documentation=http://wiki2.dovecot.org/
After=local-fs.target network.target
 
[Service]
Type=forking
User=kankurun
PIDFile=/run/kanku/kanku-triggerd.pid
ExecStart=/usr/sbin/kanku-triggerd
ExecStop=/usr/sbin/kanku-triggerd --stop

[Install]
WantedBy=multi-user.target
