# dist/systemd/kanku-worker.service
[Unit]
Description=Kanku Worker
#Documentation=man:dovecot(1)
#Documentation=http://wiki2.dovecot.org/
After=local-fs.target network.target
 
[Service]
User=kankurun
Type=forking
PIDFile=/run/kanku/kanku-worker.pid
ExecStart=/usr/sbin/kanku-worker
ExecStop=/usr/sbin/kanku-worker --stop

[Install]
WantedBy=multi-user.target
