# dist/systemd/kanku-dispatcher.service
[Unit]
Description=Kanku Dispatcher
#Documentation=man:dovecot(1)
After=local-fs.target network.target
 
[Service]
User=kankurun
Type=forking
PIDFile=/opt/kanku/var/run/kanku-dispatcher.pid
ExecStart=/opt/kanku/bin/kanku-dispatcher
ExecStop=/opt/kanku/bin/kanku-dispatcher --stop

[Install]
WantedBy=multi-user.target
