# dist/systemd/kanku-scheduler.service
[Unit]
Description=Kanku Scheduler
#Documentation=http://wiki2.dovecot.org/
After=local-fs.target network.target
 
[Service]
User=kankurun
Type=forking
PIDFile=/run/kanku/kanku-scheduler.pid
ExecStart=/usr/sbin/kanku-scheduler
ExecStop=/usr/sbin/kanku-scheduler --stop

[Install]
WantedBy=multi-user.target
