# dist/systemd/kanku-web.service
[Unit]
Description=Kanku Web
After=local-fs.target network.target
Wants=rabbitmq-server.service
 
[Service]
User=kankurun
Type=simple
ExecStart=/usr/bin/plackup /usr/lib/kanku/kanku-app.psgi
NonBlocking=yes

[Install]
WantedBy=multi-user.target
