# dist/systemd/kanku-worker.service
[Unit]
Description=Kanku Worker
After=local-fs.target network-online.target
After=rabbitmq-server.service
 
[Service]
User=kankurun
Type=forking
PIDFile=/run/kanku/kanku-worker.pid
ExecStart=/usr/sbin/kanku-worker
ExecStop=/usr/sbin/kanku-worker --stop

[Install]
WantedBy=multi-user.target
