# dist/systemd/kanku-web.service
[Unit]
Description=Kanku Web
#Documentation=man:dovecot(1)
#Documentation=http://wiki2.dovecot.org/
After=local-fs.target network.target
 
[Service]
User=kankurun
Type=simple
ExecStart=/usr/bin/plackup /usr/lib/kanku/kanku-app.psgi
NonBlocking=yes

[Install]
WantedBy=multi-user.target
